# com.tana-mcp.tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
   <dict>
      <key>Label</key>
      <string>com.tana-mcp.tunnel</string>
      <key>ProgramArguments</key>
      <array>
         <string>/opt/homebrew/bin/cloudflared</string>
         <string>tunnel</string>
         <string>run</string>
         <string>--token</string>
         <string>REPLACE_WITH_TUNNEL_TOKEN</string>
      </array>
      <key>RunAtLoad</key>
      <true />
      <key>KeepAlive</key>
      <true />
      <key>StandardOutPath</key>
      <string>/tmp/tana-mcp-tunnel.log</string>
      <key>StandardErrorPath</key>
      <string>/tmp/tana-mcp-tunnel.err</string>
   </dict>
</plist>